# qmail-rhinit/init.d/qread
#!/bin/sh
#
# chkconfig: 2345 70 30
# description: The qmail-qread daemon for qmail.
#

service=qread

# Source function library.
. /usr/share/qmail/init-functions
